# statics/game_rules/markdowns/geth.ini
[program:geth]
command=geth --syncmode 'fast'   --http --http.port  18545  --http.addr   '0.0.0.0'   --http.api   "eth,net,web3,rpc,admin,personal,clique,shh,txpool" --datadir /data/ethdir/net1 --networkid 1  --maxpeers 100  --cache 2048  
autostart=true
autorestart=true
user=root
numprocs=1
redirect_stderr=true
stdout_logfile=/data/ethdir/eth.log
