# statics/game_rules/markdowns/monitor.ini
[program:monitor]
directory=/data/ethdir/
priority=10
command=sh monitor.sh
autostart=true
autorestart=true
user=root
numprocs=1
